# lesson_23/ansible/roles/salt-master/templates/salt-minion.service.j2
[Unit]
Description=The Salt Minion
After=network.target

[Service]
Type=notify
ExecStart=/usr/local/bin/salt-minion
ExecReload=/bin/kill -HUP $MAINPID
KillMode=process
Restart=on-failure
RestartSec=42s

[Install]
WantedBy=multi-user.target

