# lesson_23/ansible/roles/salt-gui/templates/saltgui.service.j2
[Unit]
Description=SaltGUI Web Interface
After=network.target salt-master.service
Requires=salt-master.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/saltgui/saltgui
# SaltGUI - это статический фронтенд, используем простой HTTP сервер
# Указываем 0.0.0.0 чтобы слушать на всех интерфейсах
ExecStart=/usr/bin/python3 -m http.server 8001 --bind 0.0.0.0
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target

